# ssh-tunel-node-exporter.service
[Unit]
Description=Persistent SSH Tunnel Node-Exporter
After=network.target

[Service]
Restart=on-failure
RestartSec=5
User=astorf
ExecStart=/usr/bin/ssh -NTC -o ServerAliveInterval=60 -o ExitOnForwardFailure=yes -L localhost:19100:172.16.16.4:9100 -p22 monitor@remote_server_ip

[Install]
WantedBy=multi-user.target